# .config/autostart/start_apps.desktop
[Desktop Entry]
Exec=$HOME/.scripts/start_apps
Icon=dialog-scripts
Name=start apps
Path=
Type=Application
X-KDE-AutostartScript=true
